# .config/systemd/user/blackd.service
[Unit]
Description=Daemonized Black formatter for Python

[Service]
ExecStart=/usr/bin/blackd

[Install]
WantedBy=default.target
