# .config/systemd/user/gamma.service
[Unit]
Description=Apply gamma filter

[Service]
Type=oneshot
ExecStart="/usr/bin/bash" ~/.config/hypr/scripts/gamma.sh auto
