# .config/systemd/user/hyprpaper.service
[Unit]
Description=Hyprpaper

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=always

[Install]
WantedBy=default.target
