# .config/systemd/user/idleinhibit@.service
[Unit]
Description=Inhibit idle %I

[Service]
ExecStart=/home/loqusion/.bin/idleinhibit.py
