# .config/systemd/user/idleinhibit.service
[Unit]
Description=Inhibit idle

[Service]
Type=oneshot
ExecStart=/home/loqusion/.bin/idleinhibit.py
