# dot_config/systemd/user/autotrash.timer
[Unit]
Description=Empty trash

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
